# com.sargassum.collector.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.sargassum.collector</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/samueldemacedo/.local/bin/sargassum_run.sh</string>
    </array>

    <!-- Répertoire de travail : là où se trouvent sargassum_data.db et .env -->
    <key>WorkingDirectory</key>
    <string>/Users/samueldemacedo/Desktop</string>

    <!-- Toutes les 6 heures (21600 secondes) -->
    <key>StartInterval</key>
    <integer>21600</integer>

    <!-- Lancer aussi au chargement (démarrage du Mac / connexion utilisateur) -->
    <key>RunAtLoad</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/tmp/sargassum_collector.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/sargassum_collector.log</string>

    <!-- Relancer automatiquement si le process plante -->
    <key>KeepAlive</key>
    <false/>
</dict>
</plist>
